# deploy/schreibwerkstatt.service
[Unit]
Description=Schreibwerkstatt
After=network.target

[Service]
Type=simple
User=github-runner
WorkingDirectory=/opt/schreibwerkstatt
ExecStart=/usr/bin/node server.js
# always: der Admin-Backup-Restore beendet den Prozess bewusst, damit die
# Pending-DB beim Boot eingeswappt wird. `always` startet nach jedem Exit neu
# (ausser explizitem `systemctl stop`). Der Restore-Pfad nutzt zusaetzlich einen
# Exit-Code != 0, sodass er auch mit `on-failure` funktioniert, falls diese Unit
# nicht neu ausgerollt wurde.
Restart=always
RestartSec=5
Environment=PORT=3737
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
